# deployment/systemd/kast-celery.service
[Unit]
Description=KAST Web Celery Worker
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/kast-web
Environment="PATH=/opt/kast-web/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/kast-web"
EnvironmentFile=/opt/kast-web/.env

ExecStart=/opt/kast-web/venv/bin/celery -A celery_worker.celery worker \
    --loglevel=info \
    --logfile=/var/log/kast-web/celery.log

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
